# ansible/backend/java/files/backend.service.j2

[Unit]
Description=Backend
After=syslog.target

[Service]
User=opc
Environment=SPRING_CONFIG_NAME=application
Environment=SPRING_CONFIG_LOCATION={{backend_dest_directory}}/config/
ExecStart=/bin/bash {{backend_dest_directory}}/run.sh
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target